# mindaro-backend.service
[Unit]
Description=Mindaro Backend Application
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/aimbot/devhub/mindaro/working/mindaro-backend
ExecStart=/usr/bin/java -Xms512m -Xmx1536m -XX:+UseG1GC -XX:MaxGCPauseMillis=200 -XX:+UseStringDeduplication -XX:+OptimizeStringConcat -XX:+UseCompressedOops -XX:+UseCompressedClassPointers -Djava.awt.headless=true -Dfile.encoding=UTF-8 -jar build/libs/mindaro-backend-0.0.1-SNAPSHOT.jar
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mindaro-backend

# Resource limits for 2GB RAM system
LimitNOFILE=65536
LimitNPROC=4096

# Environment variables
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=JAVA_OPTS="-Xms512m -Xmx1536m -XX:+UseG1GC -XX:MaxGCPauseMillis=200"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/aimbot/devhub/mindaro/working/mindaro-backend/logs

[Install]
WantedBy=multi-user.target 